# kapjus-rag.service
[Unit]
Description=KapJus RAG Backend (FastAPI via Unix Socket)
After=network.target

# Service uses env file with actual API keys
[Service]
User=www-data
Group=www-data
WorkingDirectory=/var/www/html/kapjus.kaponline.com.br
EnvironmentFile=/var/www/html/kapjus.kaponline.com.br/.env
# Inicia o uvicorn escutando no Unix Socket para evitar overhead e portas abertas
ExecStart=/usr/bin/python3 src/python/processor.py
Restart=always
# Garante que o socket tenha as permissões corretas
ExecStartPost=/bin/sleep 2
ExecStartPost=/bin/chmod 666 /var/www/html/kapjus.kaponline.com.br/socket/kapjus.sock

[Install]
WantedBy=multi-user.target
